# config/systemd/user/mateos-battery-monitor.service
[Unit]
Description=MateOS Battery Monitor Check
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/share/mateos/bin/mateos-battery-monitor
Environment=DISPLAY=:0